# etc/systemd/socialpolitician-app-pocketbase.service
[Unit]
Description=PocketBase for app.socialpolitician.com
After=network.target

[Service]
Type=simple
User=doug
Group=doug
WorkingDirectory=/var/www/socialpolitician-app/pb_linux
# Use production data in pocketbase/pb_data; migrations via pb_linux/pb_migrations (symlink to pocketbase/pb_migrations)
ExecStart=/var/www/socialpolitician-app/pb_linux/pocketbase serve --http=127.0.0.1:8091 --dir=/var/www/socialpolitician-app/pocketbase
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
